# external/rc.local
#!/bin/sh -e
#
# rc.local
#
# This script is executed at the end of each multiuser runlevel.
# Make sure that the script will "exit 0" on success or any other
# value on error.
#
# In order to enable or disable this script just change the execution
# bits.
#
# By default this script does nothing.


systemctl stop ModemManager.service
systemctl stop ModemManager
systemctl disable ModemManager.service
systemctl disable ModemManager
echo 3 > /proc/sys/kernel/printk
sleep 3;
mknod /dev/stpwmt c 190 0
mknod /dev/stpgps c 191 0
mknod /dev/stpbt  c 192 0
mknod /dev/fm     c 193 0
mknod /dev/wmtWifi c 153 0
/usr/local/sbin/wmt_loader
/usr/local/sbin/6620_launcher -p /etc/firmware  &
sleep 1
# enable wifi
echo 1 > /dev/wmtWifi
sleep 1
#enable bt
echo 2 > /dev/wmtWifi
sleep 1
ifup wlan0 &
exit 0
